# infra/launchd/com.joel.system-bus-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joel.system-bus-worker</string>

    <key>UserName</key>
    <string>joel</string>

    <key>GroupName</key>
    <string>staff</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/joel/.local/bin/worker-supervisor</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joel/Code/joelhooks/joelclaw/infra/worker-supervisor</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/joel</string>
        <key>PATH</key>
        <string>/Users/joel/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>INNGEST_DEV</key>
        <string>0</string>
        <key>WORKER_ROLE</key>
        <string>host</string>
        <key>VAULT_PATH</key>
        <string>/Users/joel/Vault</string>
        <key>NODE_EXTRA_CA_CERTS</key>
        <string>/etc/ssl/cert.pem</string>
        <key>OTEL_STORE</key>
        <string>forward</string>
        <key>OTEL_FORWARD_URL</key>
        <string>http://flagg:3111/observability/emit</string>
        <key>OTEL_FORWARD_TIMEOUT_MS</key>
        <string>2500</string>
        <key>OTEL_TYPESENSE_PROJECTION</key>
        <string>0</string>
        <key>GATEWAY_CHANNEL_OTEL_PROBES_ENABLED</key>
        <string>0</string>
        <key>SYSTEM_HEALTH_OTEL_TYPESENSE_SIGNAL_CHECKS_ENABLED</key>
        <string>0</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>/Users/joel/.local/log/system-bus-worker.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/joel/.local/log/system-bus-worker.err</string>

    <key>ThrottleInterval</key>
    <integer>5</integer>
</dict>
</plist>
