# infra/launchd/com.joel.agent-secrets.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joel.agent-secrets</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/joel/Code/joelhooks/joelclaw/infra/agent-secrets-daemon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joel/Code/joelhooks/joelclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/joel</string>
        <key>PATH</key>
        <string>/Users/joel/.local/bin:/Users/joel/.bun/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/joel/.local/log/agent-secrets.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/joel/.local/log/agent-secrets.err</string>
</dict>
</plist>
