# infra/launchd/com.joel.content-sync-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.joel.content-sync-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/joel/Code/joelhooks/joelclaw/scripts/content-sync-watcher.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/joel/Code/joelhooks/joelclaw</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/joel</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/Users/joel/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>WatchPaths</key>
    <array>
        <string>/Users/joel/Vault/docs/decisions</string>
        <string>/Users/joel/Vault/Resources/discoveries</string>
    </array>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/tmp/content-sync-watcher.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/content-sync-watcher.log</string>
</dict>
</plist>
