# infra/central/native/com.joelclaw.central.session-index-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.joelclaw.central.session-index-health</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/Shared/joelclaw/bin/session-index-health</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>CENTRAL_ROOT</key>
    <string>/Users/Shared/joelclaw</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Users/Shared/joelclaw/logs/session-index-health.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/Shared/joelclaw/logs/session-index-health.err.log</string>
</dict>
</plist>
